# packages/pgcornac/pgCornac-2.0.0-py3-none-any.whl/cornac/files/cornac-scheduler.service
[Unit]
Description=RDS-compatible Managed-Postgres Scheduler
After=network.target

[Service]
Type=simple
User=cornac-web
Group=cornac-web
Environment=SYSTEMD=1
ExecStart=/opt/cornac/bin/cornac-shell cornac --verbose periodiq

[Install]
WantedBy=multi-user.target
