# packages/kwollect/kwollect-0.12.0.tar.gz/kwollect-0.12.0/debian/kwollector.service
[Unit]
Description=Kwollector service
After=network.target

[Service]
ExecStart=/usr/bin/kwollector /etc/kwollect/kwollector.conf

[Install]
WantedBy=multi-user.target
