# packages/raspcuterie/raspcuterie-0.1.5.tar.gz/raspcuterie-0.1.5/raspcuterie.service
[Unit]
Description=Raspcuterie dashboard
After=network.target

[Service]
User=pi
Environment=FLASK_APP=raspcuterie.app
Environment=FLASK_ENV=production
WorkingDirectory=/home/pi/raspcuterie-pi
ExecStart=/home/pi/.virtualenvs/raspcuterie/bin/flask run --host 0.0.0.0 --port 5000
Restart=always

[Install]
WantedBy=multi-user.target